# docker/files/faasd/node-exporter.service
[Unit]
Description=node-exporter
After=faasd.service

[Service]
ExecStart=/faasd/node_exporter
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
