# docker/files/dfaasagent/dfaasagent.service
[Unit]
Description=DFaaS Agent
After=haproxy.service

[Service]
EnvironmentFile=/opt/dfaasagent/dfaasagent.env
ExecStart=/opt/dfaasagent/dfaasagent
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
