# docker/files/dfaasagent/dfaasagent.service
[Unit]
Description=DFaaSAgent
After=haproxy.service

[Service]
ExecStart=/agent/dfaasagent
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target