# docker/files/dfaasagent/dfaasagent.service
[Unit]
Description=DFaaSAgent
After=haproxy.service

[Service]
ExecStart=/usr/bin/bash -c /agent/dfaasagent
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target